# tortoized.service.in
[Unit]
Description=The diricore web service
After=network.target
Wants=network-online.target

[Service]
Restart=on-failure
Type=forking
PIDFile=@TORTOIZED_PID@
ExecStart=@TORTOIZED_EXE@ start
ExecReload=@TORTOIZED_EXE@ reload
ExecStop=@TORTOIZED_EXE@ stop

[Install]
WantedBy=multi-user.target
